# wire_sentinel.service
[Unit]
Description=Wire Sentinel - Monitors for IP address changes and updates DNS and WireGuard peers.
After=network.target

[Service]
User=wire_sentinel
Group=wire_sentinel
WorkingDirectory=/opt/wire_sentinel
ExecStart=/usr/local/bin/wire_sentinel
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
