# systemd/wire_sentinel.service
[Unit]
Description=Wire Sentinel - Monitors for IP address changes and updates DNS and WireGuard peers.
After=network.target

[Service]
WorkingDirectory=/etc/wire-sentinel
ExecStart=/opt/wire-sentinel/bin/wire_sentinel
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
